# etc/etc/init/ccsclient-safe.conf
description "safely stop ccs client"

start on starting rc
task
exec /etc/init.d/ccsclient stop
